# install/docker/conf/cron/crontab
*/2 * * * * php /var/www/flightairmap/scripts/update_db.php
